# package/network/services/dpt-web-ide-server/files/dpt-web-ide-server.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Daan Pape

START=56
STOP=20

start() {
	echo "Starting DPTechnics dpt-web-ide-server"
	/usr/sbin/dpt-web-ide-server
}
